# deploy/nic-watchdog.service.j2
[Unit]
Description=NIC watchdog for RPi5 Ethernet link drops
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nic-watchdog \
  --ping-target {{ ping_target }} \
  --cooldown {{ cooldown_seconds }}s \
  --soft-max {{ soft_recover_max }}{% if metrics_addr %} \
  --metrics-addr {{ metrics_addr }}{% endif %}

Restart=always
RestartSec=5

# Security
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/sys/class/net /run
# CAP_NET_ADMIN: ip link set / neigh flush / route flush.
# Unprivileged ICMP (IPPROTO_ICMP datagram socket) is used for liveness
# probes — no CAP_NET_RAW needed; gated by net.ipv4.ping_group_range,
# which is open by default on Debian / Raspberry Pi OS.
AmbientCapabilities=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
